# supervisord.conf
[supervisord]
nodaemon=true

[program:ha-mqtt-bridge-sync]
command=/usr/local/bin/python /app.py --port %(ENV_PORT)s --oauth "%(ENV_OAUTH)s" --bridgehost %(ENV_BRIDGEHOST)s--bridgeport %(BRIDGEPORT)s
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true